# configs/keysender.socket
[Unit]
Description=Keysender Socket
PartOf=keysender.service

[Socket]
ListenStream=0.0.0.0:1338

[Install]
WantedBy=sockets.target

